# cookbooks/radar/templates/default/api/init.erb
#!/bin/bash
#
# radar-api - this script starts and stops the radar API
#
# chkconfig: - 85 15
# description: manages the radar API using supervisord
# config: /etc/radar-api/supervisord.conf
# pidfile: /var/run/radar-api-supervisord.pid
#

. /etc/rc.d/init.d/functions

RETVAL=0

start() {
  echo -n $"Starting radar-api: "
  daemon /opt/radar-api/venv/bin/python.sh /opt/radar-api/venv/bin/supervisord -c /etc/radar-api/supervisord.conf
  RETVAL=$?
  echo
  [ $RETVAL -eq 0 ] && touch /var/lock/subsys/radar-api-supervisord
}

stop() {
  echo -n $"Stopping radar-api: "
  killproc -p /var/run/radar-api-supervisord.pid
  echo
  [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/radar-api-supervisord
}

restart() {
  stop
  start
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart|force-reload|reload)
    restart
    ;;
  condrestart)
    [ -f /var/lock/subsys/radar-api-supervisord ] && restart
    ;;
  status)
    status radar-api
    RETVAL=$?
    ;;
  *)
    echo $"Usage: $0 {start|stop|status|restart|reload|force-reload|condrestart}"
    exit 1
esac

exit $RETVAL
