# cookbooks/radar/templates/default/init.erb
#!/bin/bash
#
# radar - this script starts and stops the radar API
#
# chkconfig: - 85 15
# description: instance of supervisord for the radar API
# config: /etc/radar.conf
# pidfile: /var/run/radar.pid

. /etc/rc.d/init.d/functions

RETVAL=0

start() {
  echo -n $"Starting radar: "
  daemon /home/radar/venv/bin/python.sh /home/radar/venv/bin/supervisord -c /etc/radar.conf
  RETVAL=$?
  echo
  [ $RETVAL -eq 0 ] && touch /var/lock/subsys/radar
}

stop() {
  echo -n $"Stopping radar: "
  killproc -p /var/run/radar.pid
  echo
  [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/radar
}

restart() {
  stop
  start
}

case "$1" in
  start)
    start
    ;;
  stop) 
    stop
    ;;
  restart|force-reload|reload)
    restart
    ;;
  condrestart)
    [ -f /var/lock/subsys/radar ] && restart
    ;;
  status)
    status radar
    RETVAL=$?
    ;;
  *)
    echo $"Usage: $0 {start|stop|status|restart|reload|force-reload|condrestart}"
    exit 1
esac

exit $RETVAL
